# launchd/com.nathan.vault-curator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nathan.vault-curator</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/nathan/Projects/vault-curator/scripts/daily-curate.sh</string>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>16</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/nathan/Projects/vault-curator/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/nathan/Projects/vault-curator/logs/launchd.err.log</string>
</dict>
</plist>
